# systemd/piscine-controller.service
[Unit]
Description=Piscine Controller (.NET 10 AOT)
# network-online.target garantit que la carte réseau est configurée et
# que le routage est opérationnel avant le démarrage du service.
# Sans cela, la connexion MQTT échoue avec "Network is unreachable"
# pendant plusieurs secondes au boot (observé dans les logs).
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=klm
WorkingDirectory=/home/klm/pool
ExecStart=/home/klm/pool/PiscineController
Restart=always
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=15s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=PiscineController

[Install]
WantedBy=multi-user.target
